# rossbuddy/supervisor/storescp.ini
[program:storescp]
command=/bin/storescp --log-level %(ENV_STORESCP_LOG_LEVEL)s --aetitle %(ENV_STORESCP_AE_TITLE)s --output-directory /data --promiscuous --fork %(ENV_STORESCP_PORT)s
user=rossbuddy
process_name=%(program_name)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=%(ENV_STORESCP_LOG_FILE)s
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=10

[eventlistener:slack]
command=/usr/bin/slack.py
events=PROCESS_STATE_EXITED,PROCESS_STATE_RUNNING
